# systemd/rock_jackd.service
[Unit]
Description=jackd manager
StartLimitIntervalSec=0

[Service]
Type=oneshot
#Restart=always
#RestartSec=1
User=root
#WorkingDirectory=~/Documents/leminrock/
ExecStart=export JACK_NO_AUDIO_RESERVATION=1;/usr/bin/jackd -R -P 95 -d alsa -d hw:rockchiprk3308a -r 44100 -p 256

[Install]
WantedBy=multi-user.target
